# src/systemd/cfddns.timer
[Unit]
Description=Run cfddns.service daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=cfddns.service

[Install]
WantedBy=timers.target