# systemd/container-cfddns-updater@.service.d/cf_zone_id.conf
[Service]
Environment=CF_ZONE_ID=1234567890abcdef1234567890abcdef
