# systemd/container-cfddns-updater@.service.d/cfupdater_api_params.conf
[Service]
Environment=CF_ZONE_ID=1234567890abcdef1234567890abcdef
Environment=CF_API_TOKEN_FILE=%h/secrets/cf_dns_api_token
